# iso-builder/config/includes.chroot/etc/systemd/system/disable-sleep.service
[Unit]
Description=Disable screen blanking and sleep
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disable-sleep.sh

[Install]
WantedBy=graphical.target
